# autotools/basic/main.conf
description "An example Upstart init conf file"
author      "Eric Gustafson <egustafson in launchpad>"

start on runlevel [2345]
stop on runlevel [!2345]

exec /usr/local/bin/main